# conf/local_milter.conf
[program:nexcel-milter]
command=/opt/zimbra/LocalMilter/disclaimer_milter.py
proces_name=local-milter
priority=1
redirect_stderr=true
stdout_logfile=/var/log/local-milter.log
